# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid
loglevel=info

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:laravel-octane]
command=/bin/sh -lc 'set -e; cd "${APP_DIR:-/var/www/current-backend}"; test -f artisan; test -f vendor/autoload.php; until nc -z ${DB_HOST:-db} ${DB_PORT:-3306}; do echo waiting for mysql; sleep 2; done; until nc -z ${REDIS_HOST:-redis} ${REDIS_PORT:-6379}; do echo waiting for redis; sleep 2; done; php artisan config:clear || true; php artisan cache:clear || true; php artisan route:clear || true; php artisan view:clear || true; php artisan octane:install --server=swoole || true; if [ "${RUN_MIGRATIONS:-true}" = "true" ]; then php artisan migrate --force; else echo "RUN_MIGRATIONS=false, bỏ qua migrate."; fi; exec php artisan octane:start --server=swoole --host=0.0.0.0 --port=8000 --workers=${OCTANE_WORKERS:-4} --task-workers=${OCTANE_TASK_WORKERS:-2} --max-requests=${OCTANE_MAX_REQUESTS:-1000}'
user=%(ENV_APP_USER)s
autostart=true
autorestart=true
stopsignal=QUIT
stopasgroup=true
killasgroup=true
stopwaitsecs=60
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
